# target/linux/ramips/rt305x/base-files/etc/uci-defaults/05_fix-compat-version
. /lib/functions.sh

case "$(board_name)" in
	hootoo,ht-tm02)
		uci set system.@system[0].compat_version="2.0"
		uci commit system
		;;
esac

exit 0
